# update-hack-rule.cron
*/5 * * * *     root /usr/local/sbin/update_hack_rule.sh
